# etc/backlight.rules
RUN+="/bin/chgrp arch /sys/class/backlight/intel_backlight/brightness"
RUN+="/bin/chmod g+w /sys/class/backlight/intel_backlight/brightness"
